# ansible/roles/cluster-stale-containers/files/stale-containers.service
[Unit]
Description=Stale Containers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k3s crictl rmi --prune
